# usbstick.rules
ACTION=="remove", SUBSYSTEM=="usb", ENV{PRODUCT}=="$product_id", RUN+="/bin/su $user --command='$path/stop.sh'"
